# dotfiles/system-monitor.target
[Unit]
Description=Aggregate system-monitor target
Wants=battery-monitor.service wifi-monitor.service vpn-monitor.service