# dotfiles/openrc/.config/cron/user.crontab
*/2 * * * * $HOME/repos/dotfiles/install/artixinstall/scripts/monitor-battery.sh
*/2 * * * * $HOME/repos/dotfiles/install/artixinstall/scripts/monitor-wifi.sh
